# obexd/src/obex.service.in
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=@exec_prefix@/bluetooth/obexd

[Install]
Alias=dbus-org.bluez.obex.service
