# tools/mpris-proxy.service.in
[Unit]
Description=Bluetooth mpris proxy
Documentation=man:mpris-proxy(1)

Wants=dbus.socket
After=dbus.socket dbus.service
ConditionUser=!@system

[Service]
Type=simple
ExecStart=@PKGBINDIR@/mpris-proxy

[Install]
WantedBy=default.target
