# src/bluetooth.service.in
[Unit]
Description=Bluetooth service

[Service]
Type=dbus
BusName=org.bluez
ExecStart=@prefix@/sbin/bluetoothd -n

[Install]
WantedBy=bluetooth.target
Alias=dbus-org.bluez.service
