# src/bluetooth.service.in
[Unit]
Description=Bluetooth service

[Service]
Type=dbus
BusName=org.bluez
ExecStart=@libexecdir@/bluetoothd -n
StandardOutput=null
NotifyAccess=main
#WatchdogSec=10
#Restart=on-failure

[Install]
WantedBy=bluetooth.target
Alias=dbus-org.bluez.service
